# crontab
* * * * * cd /usr/src/app/ && ./node_modules/.bin/ts-node check-vm-status.ts >> /var/log/cron.log 2>&1